# templates/82-cdrom.rules
SUBSYSTEM=="block", KERNEL=="sr0", SYMLINK+="cdrom cdrw dvd dvdrw" RUN+="{INSTALLATION}/scripts/htpcinit-dvd PLAY 3"
